# crontab
* * * * * root cd / && ruby staticize.rb >> /var/log/staticizer.log 2>&1

